# live-build/config/includes.chroot/etc/systemd/system/sushida-kiosk.service
[Unit]
Description=Sushi-da OS Kiosk
After=NetworkManager.service sushida-wifi-setup.service
Wants=NetworkManager.service sushida-wifi-setup.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=kiosk
Group=kiosk
RuntimeDirectory=sushida-os
RuntimeDirectoryMode=0750
ExecStart=/usr/local/bin/sushida-launch
Restart=always
RestartSec=3
TimeoutStopSec=10
KillMode=control-group
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/run/sushida-os

[Install]
WantedBy=multi-user.target
