# live-build/config/includes.chroot/etc/systemd/system/sushida-wifi-setup.service
[Unit]
Description=Sushi-da constrained Wi-Fi setup service
After=NetworkManager.service
Wants=NetworkManager.service
Before=sushida-kiosk.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=wifi-setup
Group=wifi-setup
RuntimeDirectory=sushida-wifi-setup
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=restart
ExecStart=/usr/local/libexec/sushida-wifi-setup
Restart=always
RestartSec=3
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=localhost
ReadWritePaths=/var/lib/sushida-config
UMask=0077
Environment=PYTHONDONTWRITEBYTECODE=1

[Install]
WantedBy=multi-user.target
